# elastic-network-interfaces.service
[Unit]
After=udev

[Service]
ExecStart=/sbin/ec2ifscan

[Install]
WantedBy=default.target
